# mageni/BUILD/supervisord.conf
[supervisord]
nodaemon=true

[program:mageni]
command=/usr/local/sbin/scanner-start

[program:mageni-db]
command=/usr/local/sbin/mageni-sqlite --listen=127.0.0.1 --port=9390 --client-watch-interval=0

[program:mageni-web]
command=/usr/local/sbin/mageni-websvr --mlisten=127.0.0.1 --mport=9390 --http-sts --gnutls-priorities=SECURE128:SECURE192:-VERS-ALL:+VERS-TLS1.2